# systemd/nipa-checks.timer
[Unit]
Description=NIPA check fetcher

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
